# Dockerfiles/kali/conf.d/burpsuite.conf
[program:burpsuite]
command=java -jar /opt/burpsuite.jar
stdout_logfile=/var/log/burpsuite.log
stderr_logfile=/var/log/burpsuite_err.log